# tools/ledsense.conf
[program:ledsense]
command=/home/pi/raspi_ledsense/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ledsense.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
